# kapitel-01-linux-grundlagen/05-prozesse-und-dienste/assets/lab-demo.service
# LabForge Workshop 5 isolated service
[Unit]
Description=LabForge Demo-Dienst
After=basic.target

[Service]
Type=simple
ExecStart=/usr/local/lib/labforge/lab-demo-runner.py
ExecStartPost=/usr/local/lib/labforge/lab-demo-event.py start ${MAINPID}
ExecStop=/usr/local/lib/labforge/lab-demo-event.py stop-request ${MAINPID}
ExecStopPost=/usr/local/lib/labforge/lab-demo-event.py stop-result ${SERVICE_RESULT} ${EXIT_CODE} ${EXIT_STATUS}
Restart=no
TimeoutStopSec=8s
